# udev-rules/01-battery.rules
SUBSYSTEM="power_supply", ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/home/manu/bin/battery_notify", OWNER="manu"
